# deploy/gcp/shiftz.service
[Unit]
Description=shiftZ unified application
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=shiftz
Group=shiftz
WorkingDirectory=/opt/shiftz
EnvironmentFile=/etc/shiftz/shiftz.env
Environment=GOMAXPROCS=1
ExecStartPre=/usr/bin/test -x /opt/shiftz/shiftz
ExecStart=/opt/shiftz/shiftz
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=15s

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
UMask=0077

MemoryHigh=600M
MemoryMax=700M
TasksMax=128
LimitNOFILE=4096

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
